# docker/s6/taskd/run
#! /bin/sh

exec taskd server --data ${TASKDDATA}
